# LDP/howto/docbook/HighQuality-Apps-HOWTO/initscript.sgml
<programlisting>
#!/bin/sh
#
# /etc/init.d/mysystem
# Subsystem file for "MySystem" server
#
# chkconfig: 2345 95 05	<co id="ex.init.chkconfig"/>
# description: MySystem server daemon
#
# processname: MySystem
# config: /etc/MySystem/mySystem.conf
# config: /etc/sysconfig/mySystem
# pidfile: /var/run/MySystem.pid

# source function library
. /etc/rc.d/init.d/functions

# pull in sysconfig settings
[ -f /etc/sysconfig/mySystem ] &amp;&amp; . /etc/sysconfig/mySystem	<co id="ex.init.sysconfig"/>

RETVAL=0
prog="MySystem"
.
.	<co id="ex.init.any.begin"/>
.

start() {	<co id="ex.init.start"/>
	echo -n $"Starting $prog:"
	.
	.	<co id="ex.init.any.start"/>
	.
	RETVAL=$?
	[ "$RETVAL" = 0 ] &amp;&amp; touch /var/lock/subsys/$prog
	echo
}

stop() {	<co id="ex.init.stop"/>
	echo -n $"Stopping $prog:"
	.
	.	<co id="ex.init.any.stop"/>
	.
	killproc $prog -TERM
	RETVAL=$?
	[ "$RETVAL" = 0 ] &amp;&amp; rm -f /var/lock/subsys/$prog
	echo
}

reload() {	<co id="ex.init.reload"/>
	echo -n $"Reloading $prog:"
	killproc $prog -HUP
	RETVAL=$?
	echo
}

case "$1" in	<co id="ex.init.call"/>
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	reload)
		reload
		;;
	condrestart)
		if [ -f /var/lock/subsys/$prog ] ; then
			stop
			# avoid race
			sleep 3
			start
		fi
		;;
	status)
		status $prog
		RETVAL=$?
		;;
	*)	<co id="ex.init.default"/>
		echo $"Usage: $0 {start|stop|restart|reload|condrestart|status}"
		RETVAL=1
esac
exit $RETVAL
</programlisting>
